# deploy/examples/zcl-c3-stopwatch-run@.service
[Unit]
Description=ZClassic23 C3 Wall-Clock Stopwatch Run (%i) — genuinely-wiped cold-start-to-tip evidence
# Templated so multiple named configurations (e.g. distinct candidate
# binaries or peers under test) can run under distinct instance names
# without duplicating this unit. %i's EnvironmentFile is optional (leading
# '-'), so the plain instance name "default" needs no override file at
# all — ZCL_BIN/ZCL_PEER/ZCL_CS_BUDGET_SECS then fall back to
# tools/scripts/c3_stopwatch_run_and_record.sh's own defaults (peer
# 127.0.0.1:39070, the dedicated zcl-stopwatch-peer.service fixture).
After=zcl-stopwatch-peer.service
# OnFailure fires ONLY when the collector could not LOCK or APPEND its
# ledger line (tools/scripts/c3_stopwatch_run_and_record.sh always exits 0
# once the append succeeds, regardless of the run's own PASS/FAIL/SKIP/
# SEAM/STALLED-NAMED verdict — that verdict is judged, not paged, via
# `make c3-stopwatch-report`). Same collect/judge split as
# zclassic23-soak-evidence.service.
OnFailure=zcl-stopwatch-onfailure.service

[Service]
Type=oneshot
WorkingDirectory=%h/github/zclassic23
EnvironmentFile=-%h/.config/zclassic23-stopwatch/%i.env
ExecStart=%h/github/zclassic23/tools/scripts/c3_stopwatch_run_and_record.sh
# 10-min stopwatch budget (ZCL_CS_BUDGET_SECS default 600) + boot/RPC
# overhead margin.
TimeoutStartSec=700
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
