# deploy/examples/zcl-netdisrupt-run@.service
[Unit]
Description=ZClassic23 Network-Disruption Recovery Stopwatch Run (%i) — PROOF B: SIGSTOP/SIGCONT upstream, time client H* recovery
# Templated the same way as zcl-c3-stopwatch-run@.service — %i names an
# optional EnvironmentFile (leading '-') carrying ZCL_ND_CLIENT_RPCPORT /
# ZCL_ND_CLIENT_DATADIR / ZCL_ND_UPSTREAM_PID_FILE (or a bare
# ZCL_ND_UPSTREAM_PID) / ZCL_ND_CUT_SECS / ZCL_ND_BUDGET_SECS for the
# already-running client + controllable upstream peer under test. This
# unit does NOT spawn either node.
After=zcl-stopwatch-peer.service
# Same collect/judge split as zcl-c3-stopwatch-run@.service: OnFailure
# fires ONLY on a lock/append failure in
# tools/scripts/netdisrupt_stopwatch_run_and_record.sh, never on a
# legitimate FAIL/SKIP/SEAM/STALLED-NAMED run verdict — `make
# netdisrupt-stopwatch-report` judges that.
OnFailure=zcl-stopwatch-onfailure.service

[Service]
Type=oneshot
WorkingDirectory=%h/github/zclassic23
EnvironmentFile=-%h/.config/zclassic23-stopwatch/%i.env
ExecStart=%h/github/zclassic23/tools/scripts/netdisrupt_stopwatch_run_and_record.sh
# Cut window (ZCL_ND_CUT_SECS default 600s) + recovery budget
# (ZCL_ND_BUDGET_SECS default 600s) + boot/RPC overhead margin.
TimeoutStartSec=1500
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
