# deploy/examples/zcl-stopwatch-peer.service
[Unit]
Description=ZClassic23 Stopwatch Fixture Peer — dedicated serving node for the C3/PROOF-B wall-clock gates
# A DEDICATED, minimal serving peer that the C3 stopwatch dials as a
# read-only P2P client (tools/scripts/cold_start_to_tip_stopwatch.sh /
# c3_stopwatch_run_and_record.sh, default peer 127.0.0.1:39070). Keeping
# it separate from the canonical zclassic23.service (port 8033) means the
# 6-hourly stopwatch run timers can never contend with, restart, or churn
# the live/canonical node's own datadir or ports.
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/state/zclassic23-stopwatch-peer/zclassic23.bin \
    -datadir=%h/.local/state/zclassic23-stopwatch-peer \
    -port=39070 -rpcport=39071 -fsport=39072 -httpsport=39073 \
    -listen -nolegacyimport -nobgvalidation -showmetrics=0 \
    -addnode=127.0.0.1:8033 -addnode=127.0.0.1:8034
Restart=always
RestartSec=15

OOMScoreAdjust=-400
MemoryHigh=6G

[Install]
WantedBy=default.target
