# deploy/examples/zclassic23-replay-canary-nightly.timer
[Unit]
Description=ZClassic23 Replay Canary Timer — nightly from-anchor

[Timer]
# 02:30 nightly with up to 30 min jitter. DISJOINT from the future
# crash-soak slot (item 7, ~04:30) and on a DISJOINT 39xxx port base
# (canary 39050/39060; crash-soak reserves 39070).
OnCalendar=*-*-* 02:30:00
RandomizedDelaySec=1800
Persistent=true
Unit=zclassic23-replay-canary-nightly.service

[Install]
WantedBy=timers.target
