# deploy/examples/zclassic23-remote-test-node.service
# Example long-running remote test node.
#
# Install as ~/.config/systemd/user/zclassic23-test.service. Keep
# operator-specific external IP and seed peers in
# ~/.config/zclassic23/remote-test.env so this unit stays reusable.

[Unit]
Description=ZClassic23 remote test node
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=-%h/.config/zclassic23/remote-test.env

ExecStart=%h/bin/zclassic23 \
    -datadir=%h/.zclassic23-test \
    -operator-lane=soak \
    -port=18033 \
    -rpcport=18233 \
    -listen \
    -nobgvalidation \
    $ZCL_TEST_EXTERNALIP_FLAG \
    $ZCL_TEST_ADDNODE_FLAGS \
    $ZCL_LANE_SNAPSHOT_LOADER_FLAG \
    -showmetrics=0

Restart=always
RestartSec=5
RestartSteps=8
RestartMaxDelaySec=600

ExecStop=/bin/sh -c 'for i in 1 2 3; do /bin/kill -TERM "$MAINPID" 2>/dev/null || exit 0; sleep 2; done'
TimeoutStopSec=300
KillMode=control-group

# Remote bootstrap and soak need room above the observed 12G RSS plateau.
# On a dedicated remote host, add a reviewed drop-in such as MemoryMax=32G.
MemoryHigh=24G
MemorySwapMax=512M
CPUWeight=30
IOWeight=30
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TasksMax=1024
OOMScoreAdjust=300

NoNewPrivileges=yes
PrivateTmp=yes

StandardOutput=append:%h/.zclassic23-test/node.log
StandardError=append:%h/.zclassic23-test/node.log

LimitCORE=infinity
Environment="ZCL_CORE_DIR=%h/.zclassic23-test/cores"

[Install]
WantedBy=default.target
