# deploy/examples/zclassic23-soak-evidence.service
[Unit]
Description=ZClassic23 Soak Evidence Collector — hourly MVP-C6 sample (read-only)
# Appends ONE JSON line per run (soak height, zclassicd height, gap,
# systemd NRestarts, ActiveEnterTimestamp, VmRSS, ok) to
# ~/.local/state/zclassic23-soak-evidence/evidence.jsonl. READ-ONLY
# against the soak node (RPC 18242) and zclassicd (RPC 8232) — it never
# restarts, stops, or writes to either. `make soak-evidence-report`
# judges the accumulated window. After= the soak node so a boot-time
# catch-up run samples a started node; deliberately NO Requires= — an
# unreachable node is recorded as an ok:false line (a hole in the
# evidence is itself evidence), never skipped.
After=zclassic23-soak.service
# OnFailure (the house pager pattern): fires ONLY when the collector
# could not APPEND its sample — unreachable nodes are NOT a unit failure.
# The AUTHORITATIVE record remains the evidence JSONL; this just makes a
# broken collector loud.
OnFailure=zclassic23-soak-evidence-onfailure.service

[Service]
Type=oneshot
ExecStart=%h/github/zclassic23/tools/scripts/soak_evidence.sh collect
# Two local RPC reads + systemctl show + /proc read — seconds, not
# minutes. Generous margin for a loaded box.
TimeoutStartSec=300
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
