# deploy/examples/zclassic23-soak.service
[Unit]
Description=ZClassic23 Three-Guarantees Soak Probe
After=zclassic23.service
Requires=zclassic23.service

[Service]
Type=oneshot
# 540s = 9 min — finishes inside the 10-min timer cadence with margin.
# Polls /api/health + getsyncdetail every 60s; asserts peer_count >= 3,
# mirror_lag <= breach, mirror_lag_breach_severity == "none", and
# magicbean_peer_count >= 1. Output goes to the journal.
ExecStart=/home/rhett/github/zclassic23/tools/scripts/soak_assert.sh 540 60

# Generous so post-restart peer rebuild (peers can take 60-120s to climb
# from 1→3) doesn't blow up the first probe. The script returns non-zero
# only at end-of-window — interim deviations are logged but don't fail
# the unit.
TimeoutStartSec=600

StandardOutput=journal
StandardError=journal
