# deploy/examples/zclassic23-self-update.timer
[Unit]
Description=Daily ZClassic23 read-only remote source/service status

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true
Unit=zclassic23-remote-status.service

[Install]
WantedBy=timers.target
