# deploy/examples/zclassic23-replay-canary-weekly.timer
[Unit]
Description=ZClassic23 Replay Canary Timer — weekly from-genesis

[Timer]
# Sunday 03:00 with up to 60 min jitter. Distinct calendar slot + 39060
# port base so a rare overlap with the nightly cannot collide.
OnCalendar=Sun *-*-* 03:00:00
RandomizedDelaySec=3600
Persistent=true
Unit=zclassic23-replay-canary-weekly.service

[Install]
WantedBy=timers.target
