# deploy/examples/zclassic23-soak-evidence.timer
[Unit]
Description=ZClassic23 Soak Evidence Timer — hourly MVP-C6 sample

[Timer]
# Hourly cadence with small jitter. The judge tolerates one missed tick
# (hole threshold = 2x cadence + 15 min slack = 8100 s); Persistent=true
# fires a catch-up sample after suspend/downtime so collector outages
# show up as ONE bounded hole, not a silent stretch of nothing.
OnCalendar=*-*-* *:05:00
RandomizedDelaySec=120
Persistent=true
Unit=zclassic23-soak-evidence.service

[Install]
WantedBy=timers.target
