# deploy/zcl23-dev.service
[Unit]
Description=ZClassic23 DEV lane (fresh-build binary, isolated datadir+ports)
After=network-online.target
Wants=network-online.target
# The DEV lane: every commit I work on gets built and deployed HERE so the
# code is actually exercised live, never rotting unrun in git. It runs the
# immutable generation selected by the atomic `current` link under
# %h/.local/lib/zclassic23-dev, transactionally activated by
# tools/dev/deploy-dev-lane.sh. It is fully isolated from the
# operator-gated live node (%h/.zclassic-c23, ports 8033/18232) and from the
# pinned-binary soak lane (%h/.zclassic-c23-soak, 8043/18242): separate datadir
# (%h/.zclassic-c23-dev) and ports (P2P 8053, RPC 18252). `make deploy` never
# touches this unit; this unit never touches the live datadir.
# A dead-on-arrival candidate may retry briefly while the activator is probing,
# but it must not loop forever. The activator quarantines a rejected generation,
# restores last-good, resets this bounded failure counter, and starts once.
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple

# No -tor: the dev lane wants forward-progress + fresh-code exercise with
# minimal moving parts; the onion surface is the live node's job.
# Peers: co-located zclassicd oracle (8034) + live zclassic23 (8033) + soak (8043)
# PLUS external mainnet peers. localhost-only peering CANNOT satisfy the
# anti-eclipse body-source floor (>=3 healthy distinct peers) — a
# localhost-only node pins at its seed whenever announce-stream holes need
# body_fetch backfill (proven on the soak lane 2026-06-12). Same external
# set as the soak unit. The dev lane owns a consensus-bound v2 starter bundle,
# so it never merges the mutable zclassicd datadir's block/index coordinates.
ExecStart=%h/.local/lib/zclassic23-dev/current/zclassic23-dev \
    -datadir=%h/.zclassic-c23-dev \
    -operator-lane=dev \
    -port=8053 \
    -rpcport=18252 \
    -listen \
    -txindex \
    -nolegacyimport \
    -addnode=127.0.0.1:8034 \
    -addnode=127.0.0.1:8033 \
    -addnode=127.0.0.1:8043 \
    -addnode=140.174.189.17 \
    -addnode=140.174.189.3 \
    -addnode=37.187.76.79 \
    -addnode=162.55.92.62 \
    -addnode=157.90.223.151 \
    -addnode=85.239.232.93 \
    $ZCL_LANE_SNAPSHOT_LOADER_FLAG \
    -showmetrics=0

# Restart=always keeps the binary's self-recovery escape hatches effective, but
# the unit-level StartLimit above bounds a dead candidate. Transactional deploy
# then restores and verifies last-good instead of leaving the bad binary active.
Restart=always
RestartSec=5
RestartSteps=8
RestartMaxDelaySec=600

# No ExecStop: standard systemd stop is one SIGTERM to the main process, then a
# TimeoutStopSec wait, then SIGKILL. The node's SIGTERM handler
# (src/main.c signal_handler) sets the shutdown flag on the first SIGTERM and
# absorbs every later one idempotently, so the old 3-pulse ExecStop added
# nothing but moving parts. Durability no longer depends on the stop signal
# cadence: shutdown WAL-checkpoints node.db and writes the verified-clean
# marker BEFORE the slow best-effort teardown (see app_shutdown_svc), so a
# SIGKILL at the TimeoutStopSec deadline can only cost the block-index boot
# cache, never the marker. (Live/soak units carry the same 3-pulse ExecStop;
# dropping it there is an owner-gated follow-up.)
TimeoutStopSec=300
KillMode=control-group

# The dev lane must never degrade the box, the live node, or the soak lane:
# low scheduling weight, hard memory cap well under the box's 93G.
MemoryHigh=18G
MemoryMax=24G
MemorySwapMax=512M
CPUQuota=150%
CPUWeight=20
IOWeight=20
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TasksMax=1024
OOMScoreAdjust=400

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.zclassic-c23-dev
PrivateTmp=yes

StandardOutput=append:%h/.zclassic-c23-dev/node.log
StandardError=append:%h/.zclassic-c23-dev/node.log

LimitCORE=infinity
Environment="ZCL_CORE_DIR=%h/.zclassic-c23-dev/cores"

[Install]
WantedBy=default.target
