# deploy/zcl23-dev.service
[Unit]
Description=ZClassic23 DEV lane (fresh-build binary, isolated datadir+ports)
After=network-online.target
Wants=network-online.target
# The DEV lane: every commit I work on gets built and deployed HERE so the
# code is actually exercised live, never rotting unrun in git. It runs a copy
# of the freshly-built binary (%h/.local/bin/zclassic23-dev), redeployed by
# tools/dev/deploy-dev-lane.sh on every build. It is fully isolated from the
# operator-gated live node (%h/.zclassic-c23, ports 8033/18232) and from the
# pinned-binary soak lane (%h/.zclassic-c23-soak, 8043/18242): separate datadir
# (%h/.zclassic-c23-dev) and ports (P2P 8053, RPC 18252). `make deploy` never
# touches this unit; this unit never touches the live datadir.
StartLimitIntervalSec=0

[Service]
Type=simple

# No -tor: the dev lane wants forward-progress + fresh-code exercise with
# minimal moving parts; the onion surface is the live node's job.
# Peers: co-located zclassicd oracle (8034) + live zclassic23 (8033) + soak (8043)
# PLUS external mainnet peers. localhost-only peering CANNOT satisfy the
# anti-eclipse body-source floor (>=3 healthy distinct peers) — a
# localhost-only node pins at its seed whenever announce-stream holes need
# body_fetch backfill (proven on the soak lane 2026-06-12). Same external
# set as the soak unit.
ExecStart=%h/.local/bin/zclassic23-dev \
    -datadir=%h/.zclassic-c23-dev \
    -operator-lane=dev \
    -port=8053 \
    -rpcport=18252 \
    -listen \
    -txindex \
    -addnode=127.0.0.1:8034 \
    -addnode=127.0.0.1:8033 \
    -addnode=127.0.0.1:8043 \
    -addnode=140.174.189.17 \
    -addnode=140.174.189.3 \
    -addnode=37.187.76.79 \
    -addnode=162.55.92.62 \
    -addnode=157.90.223.151 \
    -addnode=85.239.232.93 \
    -showmetrics=0

# Restart=always so the binary's self-recovery escape hatches (boot_crashonly
# auto-reindex, the 10-consecutive-UTXO-failure shutdown) actually fire and the
# lane stays up across a deploy/crash. Stepped backoff so a dead-on-arrival
# build loops visibly but cheaply.
Restart=always
RestartSec=5
RestartSteps=8
RestartMaxDelaySec=600

ExecStop=/bin/sh -c 'for i in 1 2 3; do /bin/kill -TERM "$MAINPID" 2>/dev/null || exit 0; sleep 2; done'
TimeoutStopSec=300
KillMode=control-group

# The dev lane must never degrade the box, the live node, or the soak lane:
# low scheduling weight, hard memory cap well under the box's 93G.
MemoryHigh=18G
MemoryMax=24G
MemorySwapMax=512M
CPUQuota=150%
CPUWeight=20
IOWeight=20
Nice=15
IOSchedulingClass=best-effort
IOSchedulingPriority=7
TasksMax=1024
OOMScoreAdjust=400

NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.zclassic-c23-dev
PrivateTmp=yes

StandardOutput=append:%h/.zclassic-c23-dev/node.log
StandardError=append:%h/.zclassic-c23-dev/node.log

LimitCORE=infinity
Environment="ZCL_CORE_DIR=%h/.zclassic-c23-dev/cores"

[Install]
WantedBy=default.target
