# deploy/zclassic23-worktree-gc.timer
[Unit]
Description=ZClassic23 Worktree + Merged-Branch Garbage Collector Timer

[Timer]
# Daily 03:35 — after zclassic23-logrotate.timer (03:00) and
# zclassic23-simnet-nightly.timer (03:10), before
# zclassic23-bundle-export.timer (04:00) and
# zclassic23-replay-canary-nightly.timer (04:30); distinct from the
# zclassic23-test-suite.timer hourly *:20:00 slot.
OnCalendar=*-*-* 03:35:00
RandomizedDelaySec=300
Persistent=true
Unit=zclassic23-worktree-gc.service

[Install]
WantedBy=timers.target
