# deploy/zclassic23-slo-probe.timer
[Unit]
Description=ZClassic23 External SLO Prober Timer — every 60s

[Timer]
# The scoreboard cadence: 60s, matching the "staying synced" SLO this
# ledger exists to prove.
#
# OnCalendar (wall-clock), NOT OnBootSec/OnUnitActiveSec (monotonic).
# The monotonic form is exactly how this ledger silently died on
# 2026-07-20 23:50: Persistent=true is honored ONLY for OnCalendar
# timers, and OnUnitActiveSec schedules relative to the service's last
# activation in the CURRENT manager runtime — so after a user-manager
# restart the timer had no anchor, never fired again, and the external
# scoreboard went dark for 3.5 days while the node sat stalled. A
# wall-clock schedule has no anchor to lose: every minute boundary is a
# firing, and Persistent=true (valid here) fires one catch-up run after
# suspend/downtime so an outage shows up as one bounded gap, not a
# silent stretch of nothing. The 5-min pager timer
# (deploy/zclassic23-slo-pager.timer) pages if this ledger ever goes
# stale >5m again, whatever the cause.
OnCalendar=*-*-* *:*:00
RandomizedDelaySec=5
Persistent=true
Unit=zclassic23-slo-probe.service

[Install]
WantedBy=timers.target
