# deploy/zclassic23-slo-probe.service
[Unit]
Description=ZClassic23 External SLO Prober — client-viewpoint uptime ledger (read-only)
# Appends 3 JSON lines per run (one per local instance: canonical rpcport
# 18232, soak rpcport 18242, dev rpcport 18252) to
# ~/.local/state/zclassic23-slo/uptime-ledger.jsonl via tools/scripts/
# node_slo_probe.sh. This is an EXTERNAL prober — it dials each instance's
# RPC port the same way any outside client would and records exactly what
# came back, never trusting node self-reports. READ-ONLY: it never
# restarts, stops, or writes to any node. After= the canonical node so a
# boot-time catch-up run samples a started node; deliberately NO Requires=
# — an unreachable instance is recorded as a reachable:false line (a hole
# in the ledger is itself the data point), never skipped.
After=zclassic23.service

[Service]
Type=oneshot
ExecStart=%h/github/zclassic23/tools/scripts/node_slo_probe.sh collect
# 4 bounded RPC calls (canonical/soak/dev/oracle), each capped at
# ZCL_SLO_RPC_TIMEOUT_SEC (default 8s) — tens of seconds worst case.
# Generous margin for a loaded box.
TimeoutStartSec=90
Nice=10
IOSchedulingClass=idle
StandardOutput=journal
StandardError=journal
