# deploy/zclassic23-test-suite.timer
[Unit]
Description=ZClassic23 Background Full Test Suite Timer

[Timer]
# Full test_parallel evidence every hour, jittered away from interactive work.
# systemd keeps the single oneshot service non-overlapping.
OnCalendar=*-*-* *:20:00
RandomizedDelaySec=300
Persistent=true
Unit=zclassic23-test-suite.service

[Install]
WantedBy=timers.target
