# deploy/zclassic23-replay-canary-weekly.timer
[Unit]
Description=ZClassic23 Replay Canary Timer — weekly from-genesis

[Timer]
# Sunday 05:30 with up to 60 min jitter. OFF-PEAK, past the nightly canary
# (04:30+jitter) and the whole 03:xx cluster (logrotate/simnet-nightly/
# coverage — coverage is itself Sun ~03:1x-03:35, so 05:30 is comfortably
# clear even on the one day they share a calendar). Distinct calendar slot
# + 39060 port base (vs the nightly's 39050) so a rare overlap cannot
# collide. ExecStart runs through replay_canary_guard.sh — same mint-fold
# skip guard as the nightly (see that .timer's comment).
OnCalendar=Sun *-*-* 05:30:00
RandomizedDelaySec=3600
Persistent=true
Unit=zclassic23-replay-canary-weekly.service

[Install]
WantedBy=timers.target
