# deploy/zclassic23-simnet-nightly.timer
[Unit]
Description=ZClassic23 Nightly Simulator Sweep Timer

[Timer]
# Wave-2 lane B2: the checked-in tools/sim/scenarios/*.scenario corpus plus
# a bounded wire_sweep seed range run ONLY under `make chaos` today, which
# `make ci` deliberately does not call (see the Makefile comment beside
# `ci:` for the measured build-cost reasoning) — so corpus regressions were
# caught only by developer discipline. This timer closes that gap nightly.
#
# 03:10 is distinct from the existing timers: zclassic23-fuzz.timer runs
# hourly (any minute, RandomizedDelaySec=300) and
# zclassic23-test-suite.timer runs daily at *:20:00 — 03:10 avoids both.
OnCalendar=*-*-* 03:10:00
RandomizedDelaySec=300
Persistent=true
Unit=zclassic23-simnet-nightly.service

[Install]
WantedBy=timers.target
