# deploy/zclassic23-test-tmp-clean.timer
[Unit]
Description=ZClassic23 test-tmp Scratch Reclaimer Timer

[Timer]
# Daily 03:45 — right after zclassic23-worktree-gc.timer (03:35), before
# zclassic23-bundle-export.timer (04:00) and
# zclassic23-replay-canary-nightly.timer (04:30); distinct from the
# zclassic23-test-suite.timer hourly *:20:00 slot.
OnCalendar=*-*-* 03:45:00
RandomizedDelaySec=300
Persistent=true
Unit=zclassic23-test-tmp-clean.service

[Install]
WantedBy=timers.target
