# deploy/zclassic23-coverage.service
[Unit]
Description=ZClassic23 Background Coverage Lane

[Service]
Type=oneshot
WorkingDirectory=%h/github/zclassic23
# Clean-skip while any *mint* user service is active; prune old lane logs
# before/after the run without changing the latest verdict on a skip.
ExecStart=%h/github/zclassic23/tools/scripts/quality_job_guard.sh coverage
TimeoutStartSec=14400

# Coverage compiles and links a separate instrumented binary. Run it as an
# idle background proof so normal development and node service stay responsive.
Environment=ZCL_QUALITY_LOG_KEEP=8
Nice=15
CPUWeight=10
IOWeight=10
IOSchedulingClass=idle
MemoryHigh=12G
TasksMax=1024

StandardOutput=journal
StandardError=journal
