# deploy/zclassic23-coverage.service
[Unit]
Description=ZClassic23 Background Coverage Lane

[Service]
Type=oneshot
WorkingDirectory=%h/github/zclassic23
ExecStart=%h/github/zclassic23/tools/scripts/background_quality_lane.sh coverage
TimeoutStartSec=14400

# Coverage compiles and links a separate instrumented binary. Run it as an
# idle background proof so normal development and node service stay responsive.
Nice=15
CPUWeight=10
IOWeight=10
IOSchedulingClass=idle
MemoryHigh=12G
TasksMax=1024

StandardOutput=journal
StandardError=journal
