# deploy/zclassic23.service
[Unit]
Description=ZClassic23 Full Node + Block Explorer
After=network-online.target
Wants=network-online.target
# Crash-only (2026-06-11): no start limit, ever. The 2026-06-11 outage
# proved the P7.6 budget (burst=10/600s) still permanently disables the
# unit: the boot integrity FATAL drained it in ~7 min and the node sat
# systemd-FAILED for 6+ hours. The binary now owns crash-loop bounding —
# boot_crashonly requests a bounded auto-reindex (max 3 per anchor) and
# pages the operator through events/health when the budget is spent.
# systemd's only job is to never stop trying; the stepped RestartSec
# backoff below keeps a genuinely-broken binary from burning CPU.
StartLimitIntervalSec=0

[Service]
# External hang watchdog (2026-07-18): a HUNG-but-alive process (the
# internal supervisor tree wedges, but the process itself never exits)
# is invisible to Restart=always above — systemd only restarts on exit,
# and a wedge is not an exit. WatchdogSec=120 + NotifyAccess close that
# gap: boot_sd_watchdog_start (config/src/boot_sd_watchdog.c) pings
# WATCHDOG=1 every WATCHDOG_USEC/2 (~60s here) ONLY while the root
# supervisor's sweep heartbeat (util/supervisor.h) is fresh AND node
# health/boot-progress looks live; a wedged supervisor tree stops the
# ping and systemd's own watchdog timer kills the unit 120s later.
# Restart=always (not on-watchdog) stays: it covers every watchdog kill
# (a watchdog timeout counts as a failing exit) AND the clean,
# intentional exits the self-recovery escape hatches rely on.
#
# Type=notify is REQUIRED with WatchdogSec, not optional: under
# Type=simple systemd considers startup complete at fork, so the
# WatchdogSec clock starts before a single line of node boot has run —
# and boot_sd_watchdog_start (the only WATCHDOG=1 source) is the LAST
# step of boot. Any boot slower than 120s (documented cold-recovery
# boots run tens of minutes) would be SIGABRT'd mid-boot and
# Restart=always would livelock it. Type=notify holds the watchdog
# clock until READY=1 (emitted by boot_sd_watchdog_start), and
# TimeoutStartSec bounds the boot itself with a cold-recovery-sized
# budget instead of the 90s default.
Type=notify
TimeoutStartSec=1800
WatchdogSec=120
NotifyAccess=main

# P5.2: operator-specific flags (externalip + addnodes) live in
# ~/.config/zclassic23/env, NOT here. The `-` prefix makes the file
# optional: a fresh clone without the env file still starts cleanly
# and falls back to DNS seeds + addrman for peer discovery. See
# deploy/zclassic23.env.example for the expected variable shape.
EnvironmentFile=-%h/.config/zclassic23/env

ExecStart=%h/zclassic23/build/bin/zclassic23 \
    -datadir=%h/.zclassic-c23 \
    -operator-lane=canonical \
    -port=8033 \
    -rpcport=18232 \
    -listen \
    -txindex \
    -tor \
    $ZCL_EXTERNALIP_FLAG \
    $ZCL_ADDNODE_FLAGS \
    -showmetrics=0
# Local operator port contract: zclassic23 owns the canonical ZClassic23 P2P
# listener on 8033. A sibling legacy zclassicd oracle must be moved off the
# canonical listener (the operator host uses 8034) and added through
# ZCL_ADDNODE_FLAGS when needed.
# Restart=always (not on-failure) so the self-recovery escape hatches —
# process_block.c (g_shutdown_requested=1 after 10 consecutive UTXO
# failures at the same height) and boot_crashonly's auto-reindex request
# (boot exits, the restart consumes the sentinel and re-derives the UTXO
# set from blocks/) — actually fire. `systemctl stop` is still honored
# as an intentional stop. Stepped backoff (systemd 254+): fast restarts
# for the auto-recovery handoff, ramping to 10 min if the binary is
# genuinely broken, so a dead-on-arrival build loops visibly but cheaply.
Restart=always
RestartSec=5
RestartSteps=8
RestartMaxDelaySec=600

# systemd's implicit stop signal is normally enough, but this node is
# too stateful to tolerate a lost or delayed stop transition. Send a
# few idempotent SIGTERM pulses explicitly through ExecStop so
# restart/stop reliably trips the binary's shutdown flag before the
# manager starts its TimeoutStopSec clock.
ExecStop=/bin/sh -c 'for i in 1 2 3; do /bin/kill -TERM "$MAINPID" 2>/dev/null || exit 0; sleep 2; done'

# Truthful verdict, systemd side. The binary itself writes the authoritative
# `shutdown-receipt.v1` into the datadir on every clean-or-forced-but-durable
# stop (see app_shutdown_svc / util/shutdown_stagewatch.h) — a from-genesis
# mint whose work SUCCEEDED but whose teardown was force-exited past the
# durability point records outcome=forced-after-durable there, NOT a failure.
# systemd's own Result can still read `timeout`/`failed` in that case, so stamp
# it next to the receipt: operators/tools compare the two and trust the
# datadir. The leading `-` keeps a stamp failure from masking the real result.
ExecStopPost=-/bin/sh -c 'printf "service_result=%s\nexit_code=%s\nexit_status=%s\n" "$SERVICE_RESULT" "$EXIT_CODE" "$EXIT_STATUS" > %h/.zclassic-c23/shutdown-systemd-result.v1'

# WAL drain + Tor teardown budget. The earlier 90s ceiling was
# breached in the 2026-05-18 session (15 restarts in 8 hours, every
# stop SIGKILL'd because WAL fsync + Tor teardown didn't complete in
# 90s on a busy node). Each SIGKILL costs peer reputation — remote
# peers flag us as flaky, throttle reconnect attempts. Bump to 300s
# so the worst-case shutdown is clean; if it ever exceeds 300s
# something is genuinely wedged and warrants investigation.
TimeoutStopSec=300

# control-group ensures the embedded Tor pthread (and any helper
# processes) die with the parent instead of lingering.
KillMode=control-group

# Memory pressure — drop OOMScoreAdjust so the kernel evicts other user
# processes first. MemoryHigh=6G soft-caps RSS without the kill-switch
# behavior of MemoryMax. -800 (lane/os-armor, 2026-07-19) rather than the
# minimum -1000: leaves room below this unit for something the operator
# deliberately wants protected even more (e.g. an OOM-killer-immune sentinel),
# while still landing far ahead of every other unprotected process on the
# host — grounded in the prior real OOM kill of this host's user manager
# (see memory feedback_linger_services_for_critical_pipelines.md). MemoryMin=2G
# is a hard reclaim floor: the kernel will NOT reclaim this cgroup's first 2G
# under memory pressure elsewhere on the host, so the node's core working set
# survives even a sibling cgroup's runaway allocation.
OOMScoreAdjust=-800
MemoryMin=2G
MemoryHigh=6G

# Sandboxing. ProtectSystem=strict makes / read-only except the explicit
# ReadWritePaths below. Tor writes its state inside vendor/tor/etc and
# the datadir; nothing else.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.zclassic-c23 %h/zclassic23/vendor/tor/etc

# PrivateTmp=yes gives the unit a unit-private /tmp — safer than
# sharing /tmp with other services. Tor currently does not use /tmp for
# sockets in our embedded build (dynhost is in-process), so this is
# safe. If Tor ever regresses to a /tmp unix socket, drop this.
PrivateTmp=yes

StandardOutput=append:%h/.zclassic-c23/node.log
StandardError=append:%h/.zclassic-c23/node.log

# P7.7: the first SIGABRT on the live node today produced no usable
# post-mortem because the kernel's default RLIMIT_CORE is 0 under
# systemd's user manager.  Lift the ceiling and hint a core-dump
# directory via env; the binary creates the dir lazily on first abort
# if ZCL_CORE_DIR is set (no ExecStartPre=mkdir to keep the unit tidy).
# The system still needs a core_pattern that writes into the per-user
# dir — e.g. as root, once per host:
#   sysctl -w kernel.core_pattern='|/usr/lib/systemd/systemd-coredump %P %u %g %s %t %c %h'
# systemd-coredump will then archive the dump under /var/lib/systemd/coredump/.
# Alternatively, a non-systemd-coredump host can set:
#   sysctl -w kernel.core_pattern='%e.%p.core'
# and change to the ZCL_CORE_DIR before dumping.
LimitCORE=infinity
Environment="ZCL_CORE_DIR=%h/.zclassic-c23/cores"

[Install]
WantedBy=default.target
