# deploy/zclassic23-standby.service
# ZClassic23 warm-standby serving replica.
#
# Purpose: a second node that follows tip continuously on its OWN datadir and
# ports, ready to be promoted to the canonical serving identity (8033/18232,
# ~/.zclassic-c23) by deploy/zclassic23-cutover.sh. Canonical has been a single
# point of failure with no failover; this unit is the always-warm understudy.
#
# It is NOT the canonical node: it runs -operator-lane=standby (an
# automated-noncanonical lane — boots a fresh datadir unattended, follows tip),
# and it sits at OOMScoreAdjust=-400 — below canonical's -800 (canonical wins
# every OOM race) but above the dev/soak/test lanes (they are evicted first).
#
# Enable (owner action; the standby datadir/ports are set in the env file):
#   cp deploy/zclassic23-standby.service ~/.config/systemd/user/
#   cp deploy/zclassic23-standby.env.example ~/.config/zclassic23/standby.env
#   $EDITOR ~/.config/zclassic23/standby.env      # STANDBY_DATADIR/PORT/RPCPORT
#   systemctl --user daemon-reload
#   systemctl --user enable --now zclassic23-standby

[Unit]
Description=ZClassic23 Warm-Standby Serving Replica (promotable understudy)
After=network-online.target
Wants=network-online.target
# Crash-only, same doctrine as the canonical unit: the binary owns crash-loop
# bounding, systemd just never stops trying. No start-limit disable.
StartLimitIntervalSec=0

[Service]
# Type=notify + WatchdogSec closes the hung-but-alive gap exactly like the
# canonical unit: boot_sd_watchdog_start pings WATCHDOG=1 only while the
# supervisor tree + boot progress look live, and Type=notify holds the
# watchdog clock until READY=1 so a slow (cold-recovery) standby boot is not
# SIGABRT'd mid-boot.
Type=notify
TimeoutStartSec=1800
WatchdogSec=120
NotifyAccess=main

# Datadir/ports come from the env file (owner-editable, kept out of the
# tracked unit). The three Environment= lines below are safe DEFAULTS so a
# fresh clone still starts on a distinct datadir/ports without an env file;
# the EnvironmentFile (loaded AFTER these, `-` = optional) overrides them.
Environment=STANDBY_DATADIR=%h/.zclassic-c23-standby
Environment=STANDBY_PORT=8053
Environment=STANDBY_RPCPORT=18272
EnvironmentFile=-%h/.config/zclassic23/standby.env
# Operator addnode seeds (same set the canonical unit uses) live alongside the
# canonical env; sourcing it here keeps the standby following the same peers.
EnvironmentFile=-%h/.config/zclassic23/env

ExecStart=%h/zclassic23/build/bin/zclassic23 \
    -datadir=${STANDBY_DATADIR} \
    -operator-lane=standby \
    -port=${STANDBY_PORT} \
    -rpcport=${STANDBY_RPCPORT} \
    -listen \
    -txindex \
    -addnode=127.0.0.1:8033 \
    -addnode=127.0.0.1:8034 \
    $ZCL_EXTERNALIP_FLAG \
    $ZCL_ADDNODE_FLAGS \
    -showmetrics=0
# No -tor: the onion serving surface is the canonical node's job; the standby's
# only job is warm forward-progress so it can be promoted. It peers with the
# local canonical (8033) + the co-located zclassicd oracle (8034) PLUS the
# external mainnet peers from ZCL_ADDNODE_FLAGS — a localhost-only peer set
# cannot satisfy the anti-eclipse body-source floor and would pin at its seed.

# Restart=always (not on-failure) so the binary's self-recovery escape hatches
# fire; `systemctl stop` is still honored as an intentional stop (the cutover
# script relies on that). Stepped backoff matches the canonical unit.
Restart=always
RestartSec=5
RestartSteps=8
RestartMaxDelaySec=600

# Idempotent SIGTERM pulses so a stop/restart reliably trips the binary's
# shutdown flag before TimeoutStopSec starts — same as canonical.
ExecStop=/bin/sh -c 'for i in 1 2 3; do /bin/kill -TERM "$MAINPID" 2>/dev/null || exit 0; sleep 2; done'
ExecStopPost=-/bin/sh -c 'printf "service_result=%s\nexit_code=%s\nexit_status=%s\n" "$SERVICE_RESULT" "$EXIT_CODE" "$EXIT_STATUS" > ${STANDBY_DATADIR}/shutdown-systemd-result.v1'

# WAL drain budget: match the canonical unit's 300s so a busy-node shutdown is
# clean and the promoted-or-demoted datadir is durable before the cutover
# script proceeds.
TimeoutStopSec=300
KillMode=control-group

# OOM armor, standby tier. OOMScoreAdjust=-400: canonical (-800) always
# out-survives the standby under host memory pressure, but the standby
# out-survives every dev/soak/test process (>= +300). MemoryMin=1G is a hard
# reclaim floor (the kernel won't reclaim the standby's first 1G — its core
# working set survives a sibling cgroup's runaway allocation). MemoryHigh=6G
# soft-caps RSS WITHOUT MemoryMax's kill-switch, and critically WITHOUT adding
# to the finite MemoryMax host budget that check_systemd_memory_budget.sh sums.
OOMScoreAdjust=-400
MemoryMin=1G
MemoryHigh=6G

# Sandboxing, same shape as canonical. ProtectSystem=strict + an explicit
# ReadWritePath for the standby datadir; nothing else is writable.
#
# IMPORTANT (systemd limitation): setting directives — ReadWritePaths,
# StandardOutput, ZCL_CORE_DIR below — do NOT expand $ENV variables; only the
# ExecStart/ExecStop command lines do. So these three MUST be the LITERAL
# default datadir (%h/.zclassic-c23-standby), and STANDBY_DATADIR defaults to
# that same path. If you point STANDBY_DATADIR somewhere else in the env file,
# you MUST also add a drop-in
# (systemctl --user edit zclassic23-standby) overriding ReadWritePaths,
# StandardOutput, StandardError and ZCL_CORE_DIR to the new path — otherwise
# the node cannot write its own datadir under ProtectSystem=strict.
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/.zclassic-c23-standby
PrivateTmp=yes

StandardOutput=append:%h/.zclassic-c23-standby/node.log
StandardError=append:%h/.zclassic-c23-standby/node.log

# Core dumps into the standby datadir (created lazily on first abort). See the
# canonical unit for the host-side core_pattern one-liner.
LimitCORE=infinity
Environment=ZCL_CORE_DIR=%h/.zclassic-c23-standby/cores

[Install]
WantedBy=default.target
