# deploy/zclassic23-coverage.timer
[Unit]
Description=ZClassic23 Background Coverage Timer

[Timer]
# Coverage is slower than fuzz smoke and less useful hourly. Keep one fresh
# weekly report without putting it on the push path.
OnCalendar=Sun *-*-* 03:15:00
RandomizedDelaySec=1800
Persistent=true
Unit=zclassic23-coverage.service

[Install]
WantedBy=timers.target
