# deploy/zclassic23-fuzz.timer
[Unit]
Description=ZClassic23 Background Fuzz Timer

[Timer]
# Four fuzzers * 15 minutes each. If a prior run is still active, systemd keeps
# the single service instance active rather than launching overlapping fuzzers.
OnCalendar=hourly
RandomizedDelaySec=300
Persistent=true
Unit=zclassic23-fuzz.service

[Install]
WantedBy=timers.target
